# systemd/pet-workflow.timer
[Unit]
Description=Pet Anime Video Daily Workflow Timer
Requires=pet-workflow.service

[Timer]
# Run every day at 9:00 AM
OnCalendar=*-*-* 09:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
